# agent/systemd/autopatch-agent.timer
[Unit]
Description=AUTO PATCH Agent Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=60s
Persistent=true

[Install]
WantedBy=timers.target
